# energymanager/rootfs/etc/s6-overlay/s6-rc.d/energymanager/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start EnergyManager service
# ==============================================================================
bashio::log.info "Starting EnergyManager..."

exec python3 /app/run.py
