# ocpp-server/rootfs/etc/s6-overlay/s6-rc.d/ocpp-server/run
#!/command/with-contenv bashio
# ==============================================================================
# OCPP Server service
# ==============================================================================

exec python3 /run.py
